# detections/wazuh/local_rules.xml
<!--
  Lab validation rules for Caldera-emulated TTPs.
  Production deployment requires environment-specific tuning and allowlists.
  Path on server: /var/ossec/etc/rules/local_rules.xml
  Incidents: INC-2026-001, INC-2026-003, INC-2026-006
-->
<group name="caldera_lab,windows,linux">

  <!-- INC-2026-001: T1197 BITS Jobs -->
  <rule id="180001" level="8" frequency="2" timeframe="120">
    <if_group>windows</if_group>
    <match>bitsadmin.exe</match>
    <match type="pcre2">(?i)/transfer.*/download</match>
    <description>Caldera-validated: Suspicious BITS download on $(win.system.computer)</description>
    <mitre>
      <id>T1197</id>
    </mitre>
    <group>attack_stage_execution,</group>
  </rule>

  <!-- 
    INC-2026-003: T1053.005 Scheduled Task persistence
    Production Tuning Note: This rule is calibrated for lab validation. In production, hardcoded strings 
    like 'ChromeUpdate' should be omitted in favor of alerts firing when schtasks.exe or Register-ScheduledTask 
    commands execute binaries/scripts from any user-writable paths (e.g., \Temp\, \AppData\, \Users\Public\, \Downloads\).
  -->
  <rule id="180002" level="10">
    <if_group>windows</if_group>
    <match type="pcre2">(?i)(schtasks\s.*/Create|Register-ScheduledTask)</match>
    <match type="pcre2">(?i)ChromeUpdate|AppData\\Local\\Temp</match>
    <description>Suspicious scheduled task with user-writable action on $(win.system.computer)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>attack_stage_persistence,</group>
  </rule>

  <!-- INC-2026-006: T1112 RDP Registry Port Modification -->
  <rule id="180003" level="8">
    <if_group>windows</if_group>
    <field name="win.eventdata.targetObject">\\REGISTRY\\MACHINE\\SYSTEM\\ControlSet001\\Control\\Terminal Server\\WinStations\\RDP-Tcp\\PortNumber</field>
    <description>Suspicious RDP registry port modification detected on $(win.system.computer)</description>
    <mitre>
      <id>T1112</id>
    </mitre>
    <group>registry_event,security_compliance</group>
  </rule>

  <!-- INC-2026-006: T1040 Network Sniffing -->
  <rule id="180004" level="9">
    <if_group>windows,linux</if_group>
    <match type="pcre2">(?i)tcpdump|tshark</match>
    <description>Unauthorized packet sniffer (tcpdump/tshark) executed on $(hostname)</description>
    <mitre>
      <id>T1040</id>
    </mitre>
    <group>process_creation,incident_response</group>
  </rule>

  <!-- Sandcat agent indicator (lab only) -->
  <rule id="180020" level="12">
    <if_group>windows,linux</if_group>
    <match type="pcre2">(?i)sandcat|caldera.*-server</match>
    <description>Caldera sandcat agent activity detected on $(hostname)</description>
    <mitre>
      <id>T1071.001</id>
    </mitre>
    <group>attack_stage_command_and_control,</group>
  </rule>

</group>
